# init/governor-switcher-ondemand.service
[Unit]
Description=Switch to ondemand

[Service]
Type=simple
ExecStart=/usr/bin/governor-switcher o

[Install]
WantedBy=default.target
