# init/governor-switcher-performance.service
[Unit]
Description=Switch to performance

[Service]
Type=simple
ExecStart=/usr/bin/governor-switcher p

[Install]
WantedBy=default.target
